# systemd/service.timer.tpl
[Unit]
Description=Run service for a scheduled time
Requires=${ service }

[Timer]
Unit=${ service }
OnCalendar=${ schedule }

[Install]
WantedBy=timers.target
